# Services/batnot.timer
[Unit]
Description=Periodically running checks for battery.

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
